# supervisor/kcptun.conf
[program:kcptun]
directory=/usr/bin
command=kcptun -c /etc/kcptun/config.json
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/kcptun.log